# image/flash-emmc.service
[Unit]
Description=Oneshot service to flash emmc

[Service]
Type=oneshot
ExecStart=/opt/emmc-boot-write.sh

[Install]
WantedBy=default.target
